# templates/default/init.erb
#!/bin/sh
# chkconfig:   <%= @runlevels %> <%= @start_priority %> <%= @stop_priority %>
# description: <%= @description %>
# processname: <%= @prog %>
#

# Source function library.
. /etc/rc.d/init.d/functions

prog="<%= @prog %>"

[ -e /etc/sysconfig/$prog ] && . /etc/sysconfig/$prog

lockfile=/var/lock/subsys/$prog

case "$1" in
    start)
        echo -n $"Starting $prog: "
        su -l <%= @username %> -c <%= @start_command %>
        retval=$?
        echo
        if [ $retval -eq 0 ]; then
          touch $lockfile
          success
        else
          failure
        end
        return $retval
        ;;
    stop)
        echo -n $"Stopping $prog: "
        su -l <%= @username %> -c <%= @stop_command %>
        retval=$?
        if [ $retval -eq 0 ]; then
          rm -f $lockfile
          success
        else
          failure
        fi
        return $retval
        ;;
    restart)
        stop
        start
        ;;
    *)
        echo $"Usage: $0 {start|stop|restart}"
        exit 2
esac
exit $?
