# templates/default/init.erb
#!/bin/sh
# chkconfig:   <%= @runlevels %> <%= @start_priority %> <%= @stop_priority %>
# description: <%= @description %>
# processname: <%= @prog %>
#

# Source function library.
. /etc/rc.d/init.d/functions

prog="<%= @prog %>"

[ -e /etc/sysconfig/$prog ] && . /etc/sysconfig/$prog

lockfile=/var/lock/subsys/$prog

case "$1" in
    start)
        echo -n $"Starting $prog: "
        su -l <%= @username %> -c <%= @start_command %> > /dev/null 2>&1
        retval=$?
        if [ $retval -eq 0 ]; then
          touch $lockfile
          success
        else
          failure
        fi
        echo
        exit $retval
        ;;
    stop)
        echo -n $"Stopping $prog: "
        su -l <%= @username %> -c <%= @stop_command %> > /dev/null 2>&1
        retval=$?
        if [ $retval -eq 0 ]; then
          rm -f $lockfile
          success
        else
          failure
        fi
        echo
        exit $retval
        ;;
    restart)
        stop
        start
        ;;
    status)
        status -p /tmp/spark-spark-org.apache.spark.deploy.master.Master-1.pid 'master' &&
            status -p /tmp/spark-spark-org.apache.spark.deploy.worker.Worker-1.pid 'worker'
        ;;
    *)
        echo $"Usage: $0 {start|stop|restart|status}"
        exit 2
esac
exit $?
